# episodic_server/episodic-server.service
[Unit]
Description=Episodic Memory HTTP Server (claude-hooks)
After=network-online.target
Wants=network-online.target
# Don't loop on failure — wait 30s, give up after 5 attempts in 5 minutes
StartLimitIntervalSec=300
StartLimitBurst=5

[Service]
Type=simple
ExecStart=/usr/bin/python3 __REPO_PATH__/episodic_server/server.py --host __HOST__ --port __PORT__
WorkingDirectory=__REPO_PATH__
Restart=on-failure
RestartSec=30

# Don't start if episodic-memory isn't installed
ExecStartPre=/bin/sh -c 'which episodic-memory >/dev/null 2>&1'

# Graceful shutdown
TimeoutStopSec=5

# Logging to journal
StandardOutput=journal
StandardError=journal
SyslogIdentifier=episodic-server

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ReadWritePaths=/root/.config/superpowers /root/.claude /var/log
PrivateTmp=true

[Install]
WantedBy=multi-user.target
