# systemd/claude-hooks-kairos-pin.service
[Unit]
Description=claude-hooks kairos-pin -- re-flip cachedGrowthBookFeatures.tengu_kairos_cron_durable to true
Documentation=https://github.com/mann1x/claude-hooks/blob/main/claude_hooks/kairos_pin.py

# Triggered by claude-hooks-kairos-pin.path. The path unit fires on
# every modification of __HOME__/.claude.json -- this oneshot service
# re-pins the durable-cron flag if Claude Code has rolled it back.

[Service]
Type=oneshot
# install.py substitutes __REPO_PATH__, __HOME__, and __USER__. Edit
# here only if you manage the unit yourself; otherwise drop overrides
# in /etc/systemd/system/claude-hooks-kairos-pin.service.d/*.conf.
User=__USER__
Environment=HOME=__HOME__
WorkingDirectory=__REPO_PATH__
ExecStart=/usr/bin/env python3 -m claude_hooks.kairos_pin --quiet
# Path-watcher debounce: systemd already coalesces rapid file changes,
# but this keeps a hard floor against accidental tight loops if the
# Python script ever started writing on every invocation.
RestartSec=2s
# Logging to journal -- see `journalctl -u claude-hooks-kairos-pin`.
StandardOutput=journal
StandardError=journal
SyslogIdentifier=claude-hooks-kairos-pin

# Security hardening -- the service only ever touches one file.
NoNewPrivileges=true
ProtectSystem=strict
ReadWritePaths=__HOME__
PrivateTmp=true

[Install]
WantedBy=multi-user.target
