# systemd/claude-hooks-rollup.service
[Unit]
Description=claude-hooks proxy rollup — ingest JSONL into stats.db
Documentation=https://github.com/mann1x/claude-hooks/blob/main/docs/PLAN-stats-sqlite.md
After=claude-hooks-proxy.service

[Service]
Type=oneshot
User=root
WorkingDirectory=/srv/dev-disk-by-label-opt/dev/claude-hooks
Environment=CLAUDE_HOOKS_REPO=/srv/dev-disk-by-label-opt/dev/claude-hooks
ExecStart=/root/anaconda3/envs/claude-hooks/bin/python scripts/proxy_rollup.py

# Reasonable hardening — rollup only reads JSONL + writes stats.db.
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/root/.claude
PrivateTmp=yes
