# systemd/claude-hooks-dashboard.service
[Unit]
Description=claude-hooks proxy dashboard — read-only stats view
Documentation=https://github.com/mann1x/claude-hooks/blob/main/docs/proxy.md
After=claude-hooks-proxy.service
Wants=claude-hooks-proxy.service

[Service]
Type=simple
User=root
WorkingDirectory=__REPO_PATH__
Environment=CLAUDE_HOOKS_REPO=__REPO_PATH__
ExecStart=__REPO_PATH__/bin/claude-hooks-dashboard
Restart=on-failure
RestartSec=5s
KillSignal=SIGTERM
TimeoutStopSec=10s
SendSIGKILL=yes

# Same hardening posture as the proxy unit. Dashboard only reads
# stats.db + ratelimit-state.json; no network out.
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=__HOME__/.claude
PrivateTmp=yes

[Install]
WantedBy=multi-user.target
