# systemd/claude-hooks-rollup.service
[Unit]
Description=claude-hooks proxy rollup — ingest JSONL into stats.db
Documentation=https://github.com/mann1x/claude-hooks/blob/main/docs/PLAN-stats-sqlite.md
After=claude-hooks-proxy.service

[Service]
Type=oneshot
User=root
WorkingDirectory=__REPO_PATH__
Environment=CLAUDE_HOOKS_REPO=__REPO_PATH__
# Use the shim so the rollup picks up the conda env's python if
# present, system python3 otherwise.
ExecStart=__REPO_PATH__/bin/claude-hooks-rollup

# Reasonable hardening — rollup only reads JSONL + writes stats.db.
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=__HOME__/.claude
PrivateTmp=yes
